# jellyfin-trilingual-movies.service
[Unit]
Description=Jellyfin Trilingual Subtitle Automation - Movies
After=network.target

[Service]
Type=simple
User=coty
WorkingDirectory=/home/coty
ExecStart=/usr/bin/python3 /home/coty/jellyfin_trilingual_subtitle_automation.py /data/Movies --watch --interval 60
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
